# postgis-backup/crontab
PATH=/usr/local/bin:/usr/bin:/bin

0 3 * * * /backup/backup.sh >> /backup/cron.log
